# files/ddsuite.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010 Jo-Philipp Wich

START=50

DDSUITE_BIN="/bin/ddsuite"

start() {
	config_load ddsuite
	$DDSUITE_BIN wifi_sta > /tmp/ddsuite.log
}

stop() {
	config_load ddsuite
}
